# mender-auto-accept/systemd/mender-auto-accept.service
[Unit]
Description=Auto-accept pending Mender devices
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/home/jehan/retina/node-infra/mender-auto-accept/.venv/bin/python /home/jehan/retina/node-infra/mender-auto-accept/auto_accept.py
EnvironmentFile=/home/jehan/retina/node-infra/mender-auto-accept/.env

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
